# rootfs/etc/supervisor/conf.d/agentbox.conf
[program:persistd]
command=/opt/agentbox/bin/persistd watch
user=root
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:code-server]
command=/opt/agentbox/services/code-server.sh
user=user
environment=HOME="/home/user",USER="user",SHELL="/bin/bash",EDITOR="code --wait",VISUAL="code --wait",GIT_EDITOR="code --wait",KUBE_EDITOR="code --wait"
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
