# rootfs/etc/supervisor/conf.d/agentbox.conf
[program:persistd]
command=/opt/agentbox/bin/persistd watch
user=root
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:code-server]
command=/opt/agentbox/code-server.sh
user=user
environment=HOME="/home/user",USER="user",LOGNAME="user",SHELL="/bin/bash"
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
