# rootfs/etc/supervisor/conf.d/composery.conf
[program:persistd]
command=/opt/persistd/bin/persistd daemon
user=root
autorestart=true
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:composery]
command=/opt/composery/code-server.sh
user=user
environment=HOME="/home/user",USER="user",LOGNAME="user",SHELL="/bin/bash"
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# A VPS has working cron; supervisor has no init to start it, so run the daemon
# in the foreground here. systemd mode gets cron from the package's own unit.
[program:cron]
command=/usr/sbin/cron -f
user=root
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
